# deploy/systemd/iron-kdcd.service
[Unit]
Description=iron-kdcd - irondirectory Kerberos 5 KDC
Documentation=https://github.com/glennswest/irondirectory
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=iron-kdcd
Group=iron-kdcd
EnvironmentFile=-/etc/iron-kdcd/iron-kdcd.conf
ExecStart=/usr/bin/iron-kdcd
Restart=always
RestartSec=1s
LimitNOFILE=65536

# Binds :88 (UDP+TCP) as a non-root user.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Hardening. No local state -- all data lives in fastetcd -- so no
# ReadWritePaths beyond what NoNewPrivileges/ProtectSystem already allow.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
